# yuuka.service.example
[Unit]
Description=Yuuka - Discord Gemini Secretary Bot
After=network.target redis-server.service

[Service]
Type=simple
# ↓ 実行ユーザーとインストール先パスをお使いの環境に合わせて変更してください
User=youruser
WorkingDirectory=/path/to/yuuka
ExecStart=/usr/local/bin/node dist/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
# YUUKA_ENCRYPTION_SECRET 等の機密は WorkingDirectory 直下の .env から自動で読み込まれます。
# .env を使わない場合は EnvironmentFile= や Environment= で直接指定してください。
# EnvironmentFile=/path/to/yuuka/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
