# yuuka-voice.service
[Unit]
Description=Yuuka Voice Client - 音声対話クライアント (voice-client)
After=network.target yuuka.service
Wants=yuuka.service

[Service]
Type=simple
User=suki
SupplementaryGroups=audio
WorkingDirectory=/home/suki/app/yuuka/voice-client
ExecStart=/usr/local/bin/node dist/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
